# metadata/envoy/systemd/globular-envoy.service
[Unit]
Description=Globular Envoy Proxy
After=network-online.target globular-xds.service
Wants=network-online.target globular-xds.service

[Service]
Type=simple
User=globular
Group=globular
WorkingDirectory={{.StateDir}}
Environment=GLOBULAR_SERVICES_DIR={{.StateDir}}/services
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -s /run/globular/envoy/envoy-bootstrap.json ] && exit 0; sleep 1; done; echo "envoy-bootstrap.json not ready"; exit 1'
ExecStart={{.Prefix}}/bin/envoy -c /run/globular/envoy/envoy-bootstrap.json --log-level info
Restart=on-failure
RestartSec=2
LimitNOFILE=524288
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
