# metadata/discovery/systemd/globular-discovery.service
[Unit]
Description=Globular discovery
After=network-online.target globular-event.service globular-rbac.service globular-resource.service
Wants=network-online.target globular-event.service globular-rbac.service globular-resource.service

[Service]
Type=simple
User=globular
Group=globular
WorkingDirectory={{.StateDir}}/discovery
Environment=GLOBULAR_SERVICES_DIR={{.StateDir}}/services

ExecStartPre=+/bin/sh -c 'mkdir -p {{.StateDir}}/discovery && chown globular:globular {{.StateDir}}/discovery'
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -f /var/lib/globular/pki/issued/services/service.crt ] && exit 0; sleep 1; done; echo "TLS cert not ready"; exit 1'
ExecStart={{.Prefix}}/bin/discovery_server
Restart=on-failure
RestartSec=2
LimitNOFILE=524288

[Install]
WantedBy=multi-user.target
